# docker/supervisord.conf
# Configuration Supervisord pour gérer Gunicorn et Nginx

[supervisord]
nodaemon=true
logfile=/var/log/app/supervisord.log
pidfile=/var/run/supervisord.pid

# Django application avec Gunicorn
[program:gunicorn]
directory=/app
command=gunicorn --workers 4 --worker-class sync --bind 127.0.0.1:8000 --timeout 60 --access-logfile /var/log/app/gunicorn_access.log --error-logfile /var/log/app/gunicorn_error.log config.wsgi:application
user=www-data
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/var/log/app/gunicorn.log
stderr_logfile=/var/log/app/gunicorn_error.log

# Nginx web server
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/app/nginx.log
stderr_logfile=/var/log/app/nginx_error.log
user=root

# Log des erreurs
[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
